# systemd/bitfinexdumper-compress-databases.timer
[Unit]
Description=Compress bitfinexdumper's database every weeks

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target